# systemd/restic-backup.timer
[Unit]
Description=Timer to restic backup service

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
