# contrib/restic-backup.service.d/override.conf
[Service]
# Use the absolute path to directory with scripts
Environment="RESTIC_BACKUP_DIRECTORY=/home/my_user/restic-backup"